# freemodel-proxy.service
[Unit]
Description=freemodel-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/workspace/freemodel-proxy
Environment=NODE_ENV=production
ExecStart=/bin/bash -lc 'source /root/.nvm/nvm.sh && (nvm use 20 || nvm use --lts) && exec node src/index.js'
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
